# files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

pgrep -f 'socat -ly' &>/dev/null || /etc/socat.sh
[ -e /storage ] || mkdir -p /storage
mount himawari:/storage /storage -o noatime

exit 0
